# config/crontab
0	0	*/10	*	*	acme-client -v esote.net && rcctl reload httpd
